# egm-streamer@.service
[Unit]
Description=EGM State Detector for %i
After=network.target

[Service]
Type=simple
User=user
Group=user
# Adjust path to your venv or global python
# Adjust path to your venv or global python
ExecStart=/opt/egm-streamer/venv/bin/egm-streamer serve --config /etc/egm-streamer/%i.yaml
Restart=always
RestartSec=5

# Create runtime dir for status files if needed
RuntimeDirectory=egm-streamer/%i
RuntimeDirectoryMode=0755
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
